# src/supervisord.conf
[supervisord]
nodaemon=true

[program:postgresql]
command=service postgresql start

[program:master]
command=/var/src/rsaworkers/dist/rsaworkers/rsaworker