# src/supervisord.conf
[supervisord]
nodaemon=true

#[program:postgres]
#command=service postgresql start
#priority=1

[program:rsaworker]
command=/var/src/rsaworkers/dist/rsaworkers/rsaworker
priority=2
